# debian/watchdog-mux.service
[Unit]
Description=Proxmox VE watchdog multiplexer

[Service]
ExecStart=/usr/sbin/watchdog-mux
OOMScoreAdjust=-1000