# debian/watchdog-mux.socket
[Socket]
ListenStream=/run/watchdog-mux.sock

[Install]
WantedBy=sockets.target
