# debian/watchdog-mux.service
[Service]
ExecStart=/usr/sbin/watchdog-mux
